# systemd/koalabyte-menu.service
[Unit]
Description=KoalaByte Blue Boot Launcher and Jungle Menu
After=network-online.target bluetooth.service systemd-user-sessions.service
Wants=network-online.target bluetooth.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/opt/KoalaByte-Blue
Environment=PYTHONPATH=/opt/KoalaByte-Blue/pi-companion:/opt/KoalaByte-Blue
Environment=PYTHON_BIN=/opt/KoalaByte-Blue/pi-companion/.venv/bin/python
Environment=KOALABYTE_MENU_SYNC=auto
Environment=KOALABYTE_TTS=1
Environment=MENU_GRAPHICAL=1
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/test -f /opt/KoalaByte-Blue/scripts/koalabyte_blue_boot.sh
ExecStart=/usr/bin/bash /opt/KoalaByte-Blue/scripts/koalabyte_blue_boot.sh
Restart=always
RestartSec=5
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
