# systemd/koalabyte-thats-not-a-knife.service.in
[Unit]
Description=KoalaByte Blue that’s not a knife BLE guard
After=bluetooth.service multi-user.target
Wants=bluetooth.service

[Service]
Type=simple
User={{SERVICE_USER}}
WorkingDirectory={{REPO_ROOT}}
Environment=PYTHONPATH={{REPO_ROOT}}/pi-companion
ExecStart={{PYTHON_BIN}} {{REPO_ROOT}}/scripts/run_thats_not_a_knife_loop.py --interval {{INTERVAL_SECONDS}} --threshold {{THRESHOLD}} --xp-cooldown {{XP_COOLDOWN_SECONDS}} --output-dir {{REPO_ROOT}}/logs/thats_not_a_knife --state-path {{REPO_ROOT}}/logs/thats_not_a_knife/guard_state.json --block-path {{REPO_ROOT}}/logs/thats_not_a_knife/ble_workflow_block.json --xp-path {{REPO_ROOT}}/logs/killerkoala/xp_state.json --log-glob {{REPO_ROOT}}/logs/koala_bluez/*.json --log-glob {{REPO_ROOT}}/logs/killerkoala_voice/*.json --log-glob /var/log/syslog --log-glob /var/log/messages --log-glob /var/log/kern.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
