# systemd/koalabyte-doctor.service
[Unit]
Description=KoalaByte Blue Doctor Check
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=oneshot
WorkingDirectory=/opt/KoalaByte-Blue
Environment=PYTHONPATH=/opt/KoalaByte-Blue/pi-companion:/opt/KoalaByte-Blue
ExecStart=/usr/bin/bash /opt/KoalaByte-Blue/scripts/koalabyte_doctor.sh --quick
RemainAfterExit=yes
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
